# analytics-web/deploy/lavozsalsa-analytics.service
[Unit]
Description=La Voz Salsa Analytics Web
After=network.target

[Service]
Type=simple
User=lavozsalsa
Group=lavozsalsa
WorkingDirectory=/opt/lavozsalsa-analytics/current/analytics-web
EnvironmentFile=/opt/lavozsalsa-analytics/shared/.env
ExecStart=/usr/bin/node --no-warnings /opt/lavozsalsa-analytics/current/analytics-web/src/server.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
